# salt/states/hashicorp/files/nomad.service
[Unit]
Description=Nomad Agent
Documentation=https://www.nomadproject.io/docs/
Wants=network-online.target
After=network-online.target

# If running with Consul integration
ConditionFileNotEmpty=/etc/nomad.d/nomad.hcl
Wants=consul.service
After=consul.service

# If running with Vault integration
After=vault.service

[Service]
Type=notify
User=nomad
Group=nomad
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/nomad agent -config /etc/nomad.d
KillMode=process
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitNPROC=infinity
TasksMax=infinity

# Security hardening
NoNewPrivileges=false  # Set to false because Nomad needs to spawn containers
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nomad /var/log/nomad /opt/nomad
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictNamespaces=false  # Nomad needs namespace access for containers
RestrictSUIDSGID=true
LockPersonality=true

# Environment configuration
Environment="GOMAXPROCS=2"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nomad

# OOM handling
OOMScoreAdjust=-1000

[Install]
WantedBy=multi-user.target