# assets/services/delphi-agent-enricher.service
# /etc/systemd/system/delphi-agent-enricher.service
[Unit]
Description=Delphi Emailer Service
After=network.target

[Service]
Type=simple
User=stanley
ExecStart=/usr/bin/python3 /usr/local/bin/delphi-agent-enricher.py
WorkingDirectory=/home/ubuntu
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
MemoryMax=1G
CPUQuota=80%
TasksMax=100
TimeoutStartSec=90
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/stackstorm/packs/delphi/.env

[Install]
WantedBy=multi-user.target