# assets/services/email-structurer.service
[Unit]
Description=Delphi Email Structurer Worker
Documentation=https://wiki.cybermonkey.net.au/delphi/email-structurer
After=network.target postgresql.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=simple
User=stanley
Group=stanley
WorkingDirectory=/opt/stackstorm/packs/delphi

# Main executable
ExecStart=/usr/local/bin/email-structurer.py

# Environment configuration
Environment="PYTHONPATH=/usr/local/lib/python3.11/site-packages:/opt/stackstorm/packs/delphi"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/stackstorm/packs/delphi/.env

# Email structurer specific configuration
Environment="DELPHI_PARSER_TYPE=standard"
Environment="DELPHI_CUSTOM_SECTIONS="

# Logging configuration
Environment="LOG_FILE=/var/log/stackstorm/email-structurer.log"
Environment="MAX_LOG_SIZE=5242880"

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=8192
LimitNPROC=256

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stackstorm /opt/stackstorm/packs/delphi /tmp

# Monitoring and health checks
TimeoutStartSec=30
TimeoutStopSec=15
WatchdogSec=60

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-structurer

[Install]
WantedBy=multi-user.target
Alias=delphi-email-structurer.service

# Service management commands:
# sudo systemctl enable email-structurer.service
# sudo systemctl start email-structurer.service  
# sudo systemctl status email-structurer.service
# sudo systemctl logs -f email-structurer.service
# sudo systemctl reload email-structurer.service