# assets/services/email-sender.service
[Unit]
Description=Delphi Email Sender Worker
Documentation=https://wiki.cybermonkey.net.au/delphi/email-sender
After=network.target postgresql.service email-formatter.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=simple
User=stanley
Group=stanley
WorkingDirectory=/opt/stackstorm/packs/delphi

# Main executable
ExecStart=/usr/local/bin/email-sender.py

# Environment configuration
Environment="PYTHONPATH=/usr/local/lib/python3.11/site-packages:/opt/stackstorm/packs/delphi"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/stackstorm/packs/delphi/.env

# Email sender specific configuration
Environment="MAX_RECIPIENTS=10"
Environment="EMAIL_RETRY_ATTEMPTS=3"
Environment="EMAIL_RETRY_DELAY=60"

# SMTP configuration (should be in .env file)
# SMTP_HOST=smtp.example.com
# SMTP_PORT=587
# SMTP_USER=user@example.com
# SMTP_PASS=password
# SMTP_TLS=true
# FROM_EMAIL=delphi@cybermonkey.net.au
# DEFAULT_TO_EMAIL=security@cybermonkey.net.au

# Logging configuration
Environment="LOG_FILE=/var/log/stackstorm/email-sender.log"
Environment="MAX_LOG_SIZE=5242880"

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=8192
LimitNPROC=256

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stackstorm /opt/stackstorm/packs/delphi /tmp

# Monitoring and health checks
TimeoutStartSec=30
TimeoutStopSec=15
WatchdogSec=60

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-sender

[Install]
WantedBy=multi-user.target
Alias=delphi-email-sender.service

# Service management commands:
# sudo systemctl enable email-sender.service
# sudo systemctl start email-sender.service  
# sudo systemctl status email-sender.service
# sudo systemctl logs -f email-sender.service
# sudo systemctl reload email-sender.service