# assets/services/llm-worker.service
# /etc/systemd/system/llm-worker.service
[Unit]
Description=LLM Worker for Summarizing Security Alerts
After=network.target postgresql.service

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/llm-worker.py
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
MemoryMax=1G
CPUQuota=80%
TasksMax=100
TimeoutStartSec=90
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
Type=simple
User=stanley
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/stackstorm/packs/delphi/.env
WorkingDirectory=/usr/local/bin

[Install]
WantedBy=multi-user.target