# assets/services/parser-monitor.service
[Unit]
Description=Wazuh Parser Health Monitor
Documentation=https://wiki.cybermonkey.net.au/wazuh/parser-monitor
After=network.target postgresql.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=oneshot
User=stanley
Group=stanley
WorkingDirectory=/opt/stackstorm/packs/wazuh

# Main executable (runs once on demand)
ExecStart=/usr/local/bin/parser-monitor.py --health

# Environment configuration
Environment="PYTHONPATH=/usr/local/lib/python3.11/site-packages:/opt/stackstorm/packs/wazuh"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/stackstorm/packs/wazuh/.env

# Parser monitor specific configuration
Environment="MONITOR_DASHBOARD_MODE=health"
Environment="MONITOR_OUTPUT_FORMAT=text"

# Logging configuration
Environment="LOG_FILE=/var/log/stackstorm/parser-monitor.log"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stackstorm /opt/stackstorm/packs/wazuh /tmp

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parser-monitor

[Install]
WantedBy=multi-user.target
Alias=wazuh-parser-monitor.service

# Service management commands:
# sudo systemctl enable parser-monitor.service  
# sudo systemctl start parser-monitor.service
# sudo systemctl status parser-monitor.service
# 
# Manual monitoring commands:
# /usr/local/bin/parser-monitor.py                    # Full dashboard
# /usr/local/bin/parser-monitor.py --health           # Health summary
# /usr/local/bin/parser-monitor.py --continuous       # Live monitoring
# /usr/local/bin/parser-monitor.py --failures         # Recent failures
# /usr/local/bin/parser-monitor.py --performance      # Performance metrics