# assets/services/prompt-ab-tester.service
[Unit]
Description=Delphi Prompt A/B Testing Worker
Documentation=https://wiki.cybermonkey.net.au/delphi/ab-testing
After=network.target postgresql.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=simple
User=stanley
Group=stanley
WorkingDirectory=/opt/stackstorm/packs/delphi

# Main executable
ExecStart=/usr/local/bin/prompt-ab-tester.py

# Environment configuration
Environment="PYTHONPATH=/usr/local/lib/python3.11/site-packages:/opt/stackstorm/packs/delphi"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/stackstorm/packs/delphi/.env

# A/B Testing specific configuration
Environment="EXPERIMENT_CONFIG_FILE=/opt/delphi/ab-test-config.json"
Environment="SYSTEM_PROMPTS_DIR=/opt/eos/assets/system-prompts"
Environment="DEFAULT_PROMPT_FILE=/opt/system-prompt.txt"

# Logging configuration
Environment="LOG_FILE=/var/log/stackstorm/prompt-ab-tester.log"
Environment="HEARTBEAT_FILE=/var/log/stackstorm/prompt-ab-tester.heartbeat"
Environment="METRICS_FILE=/var/log/stackstorm/ab-test-metrics.log"
Environment="MAX_LOG_SIZE=10485760"

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=8192
LimitNPROC=256

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stackstorm /opt/delphi /tmp

# Monitoring and health checks
TimeoutStartSec=30
TimeoutStopSec=15
WatchdogSec=60

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prompt-ab-tester

[Install]
WantedBy=multi-user.target
Alias=delphi-ab-tester.service

# Service management commands:
# sudo systemctl enable prompt-ab-tester.service
# sudo systemctl start prompt-ab-tester.service  
# sudo systemctl status prompt-ab-tester.service
# sudo systemctl logs -f prompt-ab-tester.service
# sudo systemctl reload prompt-ab-tester.service