# salt/api/eos-salt-api.service
[Unit]
Description=EOS Salt API Server
After=network.target salt-master.service
Requires=salt-master.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/eos/salt/api
ExecStart=/usr/bin/python3 /opt/eos/salt/api/cluster_api.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/eos/salt/api

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eos-salt-api

# Security
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target